# system_files/shared/etc/systemd/system/plocate-updatedb.service.d/override.conf
[Service]
IOSchedulingClass=Idle
IOSchedulingPriority=7
CPUSchedulingClass=Idle
